# supervisord.conf
[unix_http_server]
file=/tmp/alphadog_supervisor.sock

[supervisord]
logfile=/tmp/alphadog_supervisord.log
logfile_maxbytes=10MB
logfile_backups=3
loglevel=info
pidfile=/tmp/alphadog_supervisord.pid
nodaemon=false
user=wuqisheng

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[supervisorctl]
serverurl=unix:///tmp/alphadog_supervisor.sock

[program:pipeline]
command=python3 pipeline/main.py
directory=/Users/wuqisheng/.openclaw/workspace/alphadog
environment=BINANCE_TESTNET="true"
autostart=true
autorestart=true
startretries=3
stderr_logfile=/tmp/alphadog_pipeline_err.log
stdout_logfile=/tmp/alphadog_pipeline.log
stdout_logfile_maxbytes=10MB
stdout_logfile_backups=3

[program:engine]
command=python3 engine/run.py
directory=/Users/wuqisheng/.openclaw/workspace/alphadog
environment=BINANCE_TESTNET="true"
autostart=true
autorestart=true
startretries=3
stderr_logfile=/tmp/alphadog_engine_err.log
stdout_logfile=/tmp/alphadog_engine.log
stdout_logfile_maxbytes=10MB
stdout_logfile_backups=3

[program:api]
command=uvicorn api.main:app --host 0.0.0.0 --port 8000
directory=/Users/wuqisheng/.openclaw/workspace/alphadog
environment=BINANCE_TESTNET="true"
autostart=true
autorestart=true
startretries=3
stderr_logfile=/tmp/alphadog_api_err.log
stdout_logfile=/tmp/alphadog_api.log
stdout_logfile_maxbytes=10MB
stdout_logfile_backups=3

[program:trading]
command=python3 trader/runner.py
 directory=/Users/wuqisheng/.openclaw/workspace/alphadog
environment=BINANCE_TESTNET="true"
autostart=true
autorestart=true
startretries=3
stderr_logfile=/tmp/alphadog_trading_err.log
stdout_logfile=/tmp/alphadog_trading.log
stdout_logfile_maxbytes=10MB
stdout_logfile_backups=3

[program:frontend]
command=npx vite --host 0.0.0.0 --port 3000
directory=/Users/wuqisheng/.openclaw/workspace/alphadog/frontend
autostart=true
autorestart=true
startretries=3
stderr_logfile=/tmp/alphadog_frontend_err.log
stdout_logfile=/tmp/alphadog_frontend.log
stdout_logfile_maxbytes=10MB
stdout_logfile_backups=3
